# packaging/aur/attackshark-battery-bridge.install
post_install() {
  cat <<'EOF'
attackshark-battery-bridge was installed.

Default config:
  /etc/attackshark-battery-bridge/config.toml

Enable and start the daemon:
  sudo systemctl enable --now attackshark-battery-bridge.service
EOF
}

post_upgrade() {
  post_install
}
